# main.cron
00 23 * * * python3 /home/deaise/bopo/project/main.py
